# ports/opt/bluez/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/bluetooth/
-rw-r--r--	root/root	etc/bluetooth/audio.conf
-rw-r--r--	root/root	etc/bluetooth/input.conf
-rw-r--r--	root/root	etc/bluetooth/main.conf
-rw-r--r--	root/root	etc/bluetooth/network.conf
-rw-r--r--	root/root	etc/bluetooth/rfcomm.conf
-rw-r--r--	root/root	etc/bluetooth/serial.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/bluetoothd
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/udev/
-rwxr-xr-x	root/root	lib/udev/hid2hci
drwxr-xr-x	root/root	lib/udev/rules.d/
-rw-r--r--	root/root	lib/udev/rules.d/97-bluetooth-hid2hci.rules
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/bluez-avinfo
-rwxr-xr-x	root/root	usr/bin/bluez-dfubabel
-rwxr-xr-x	root/root	usr/bin/bluez-hcieventmask
-rwxr-xr-x	root/root	usr/bin/bluez-hcisecfilter
-rwxr-xr-x	root/root	usr/bin/bluez-hsmicro
-rwxr-xr-x	root/root	usr/bin/bluez-hsplay
-rwxr-xr-x	root/root	usr/bin/bluez-hsplaye
-rwxr-xr-x	root/root	usr/bin/bluez-list-devices
-rwxr-xr-x	root/root	usr/bin/bluez-monitor-bluetooth
-rwxr-xr-x	root/root	usr/bin/bluez-ppporc
-rwxr-xr-x	root/root	usr/bin/bluez-simple-agent
-rwxr-xr-x	root/root	usr/bin/bluez-simple-agente
-rwxr-xr-x	root/root	usr/bin/bluez-simple-endpoint
-rwxr-xr-x	root/root	usr/bin/bluez-simple-player
-rwxr-xr-x	root/root	usr/bin/bluez-simple-service
-rwxr-xr-x	root/root	usr/bin/bluez-test-adapter
-rwxr-xr-x	root/root	usr/bin/bluez-test-attrib
-rwxr-xr-x	root/root	usr/bin/bluez-test-audio
-rwxr-xr-x	root/root	usr/bin/bluez-test-device
-rwxr-xr-x	root/root	usr/bin/bluez-test-devicee
-rwxr-xr-x	root/root	usr/bin/bluez-test-discovery
-rwxr-xr-x	root/root	usr/bin/bluez-test-discoverye
-rwxr-xr-x	root/root	usr/bin/bluez-test-health
-rwxr-xr-x	root/root	usr/bin/bluez-test-health-sink
-rwxr-xr-x	root/root	usr/bin/bluez-test-input
-rwxr-xr-x	root/root	usr/bin/bluez-test-manager
-rwxr-xr-x	root/root	usr/bin/bluez-test-managere
-rwxr-xr-x	root/root	usr/bin/bluez-test-nap
-rwxr-xr-x	root/root	usr/bin/bluez-test-network
-rwxr-xr-x	root/root	usr/bin/bluez-test-oob
-rwxr-xr-x	root/root	usr/bin/bluez-test-proximity
-rwxr-xr-x	root/root	usr/bin/bluez-test-sap-server
-rwxr-xr-x	root/root	usr/bin/bluez-test-serial
-rwxr-xr-x	root/root	usr/bin/bluez-test-serial-proxy
-rwxr-xr-x	root/root	usr/bin/bluez-test-service
-rwxr-xr-x	root/root	usr/bin/bluez-test-telephony
-rwxr-xr-x	root/root	usr/bin/bluez-test-thermometer
-rwxr-xr-x	root/root	usr/bin/ciptool
-rwxr-xr-x	root/root	usr/bin/dfutool
-rwxr-xr-x	root/root	usr/bin/dund
-rwxr-xr-x	root/root	usr/bin/gatttool
-rwxr-xr-x	root/root	usr/bin/hcitool
-rwxr-xr-x	root/root	usr/bin/hidd
-rwxr-xr-x	root/root	usr/bin/l2ping
-rwxr-xr-x	root/root	usr/bin/pand
-rwxr-xr-x	root/root	usr/bin/rfcomm
-rwxr-xr-x	root/root	usr/bin/sdptool
drwxr-xr-x	root/root	usr/etc/
drwxr-xr-x	root/root	usr/etc/dbus-1/
drwxr-xr-x	root/root	usr/etc/dbus-1/system.d/
-rw-r--r--	root/root	usr/etc/dbus-1/system.d/bluetooth.conf
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/bluetooth/
-rw-r--r--	root/root	usr/include/bluetooth/a2mp.h
-rw-r--r--	root/root	usr/include/bluetooth/bluetooth.h
-rw-r--r--	root/root	usr/include/bluetooth/bnep.h
-rw-r--r--	root/root	usr/include/bluetooth/cmtp.h
-rw-r--r--	root/root	usr/include/bluetooth/hci.h
-rw-r--r--	root/root	usr/include/bluetooth/hci_lib.h
-rw-r--r--	root/root	usr/include/bluetooth/hidp.h
-rw-r--r--	root/root	usr/include/bluetooth/l2cap.h
-rw-r--r--	root/root	usr/include/bluetooth/mgmt.h
-rw-r--r--	root/root	usr/include/bluetooth/rfcomm.h
-rw-r--r--	root/root	usr/include/bluetooth/sco.h
-rw-r--r--	root/root	usr/include/bluetooth/sdp.h
-rw-r--r--	root/root	usr/include/bluetooth/sdp_lib.h
-rw-r--r--	root/root	usr/include/bluetooth/uuid.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/alsa-lib/
-rwxr-xr-x	root/root	usr/lib/alsa-lib/libasound_module_ctl_bluetooth.la
-rwxr-xr-x	root/root	usr/lib/alsa-lib/libasound_module_ctl_bluetooth.so
-rwxr-xr-x	root/root	usr/lib/alsa-lib/libasound_module_pcm_bluetooth.la
-rwxr-xr-x	root/root	usr/lib/alsa-lib/libasound_module_pcm_bluetooth.so
drwxr-xr-x	root/root	usr/lib/bluetooth/
drwxr-xr-x	root/root	usr/lib/bluetooth/plugins/
drwxr-xr-x	root/root	usr/lib/cups/
drwxr-xr-x	root/root	usr/lib/cups/backend/
-rwxr-xr-x	root/root	usr/lib/cups/backend/bluetooth
-rwxr-xr-x	root/root	usr/lib/libbluetooth.la
lrwxrwxrwx	root/root	usr/lib/libbluetooth.so -> libbluetooth.so.3.13.0
lrwxrwxrwx	root/root	usr/lib/libbluetooth.so.3 -> libbluetooth.so.3.13.0
-rwxr-xr-x	root/root	usr/lib/libbluetooth.so.3.13.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/bluez.pc
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/bccmd
-rwxr-xr-x	root/root	usr/sbin/bluetoothd
-rwxr-xr-x	root/root	usr/sbin/bluez-avctrl
-rwxr-xr-x	root/root	usr/sbin/bluez-hid2hci
-rwxr-xr-x	root/root	usr/sbin/hciattach
-rwxr-xr-x	root/root	usr/sbin/hciconfig
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/alsa/
-rw-r--r--	root/root	usr/share/alsa/bluetooth.conf
drwxr-xr-x	root/root	usr/share/dbus-1/
drwxr-xr-x	root/root	usr/share/dbus-1/system-services/
-rw-r--r--	root/root	usr/share/dbus-1/system-services/org.bluez.service
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/bluez-dfubabel.1.gz
-rw-r--r--	root/root	usr/share/man/man1/ciptool.1.gz
-rw-r--r--	root/root	usr/share/man/man1/dfutool.1.gz
-rw-r--r--	root/root	usr/share/man/man1/dund.1.gz
-rw-r--r--	root/root	usr/share/man/man1/hcitool.1.gz
-rw-r--r--	root/root	usr/share/man/man1/hidd.1.gz
-rw-r--r--	root/root	usr/share/man/man1/pand.1.gz
-rw-r--r--	root/root	usr/share/man/man1/rfcomm.1.gz
-rw-r--r--	root/root	usr/share/man/man1/sdptool.1.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/bccmd.8.gz
-rw-r--r--	root/root	usr/share/man/man8/bluetoothd.8.gz
-rw-r--r--	root/root	usr/share/man/man8/bluez-avctrl.8.gz
-rw-r--r--	root/root	usr/share/man/man8/bluez-hid2hci.8.gz
-rw-r--r--	root/root	usr/share/man/man8/hciattach.8.gz
-rw-r--r--	root/root	usr/share/man/man8/hciconfig.8.gz
-rw-r--r--	root/root	usr/share/man/man8/hid2hci.8.gz
-rw-r--r--	root/root	usr/share/man/man8/l2ping.8.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/bluetooth/
